# agent/systemd/wardenhub-agent.timer
[Unit]
Description=WardHub Security Audit Timer
After=network-online.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=3h
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
